# crons/frontend_data.cron
##External data fetch
0 6,18 * * * root /usr/bin/flock -n /tmp/update_frontend_yahoo_estimates.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_estimates.php 2>&1 | tail -n 200 | mail -s "Yahoo Estimates Cron Output" cronreports.i9xgk@zapiermail.com'
0 17,19,21,23 * * * root /usr/bin/flock -n /tmp/update_frontend_yahoo_daily.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_daily.php 2>&1 | tail -n 200 | mail -s "Yahoo Daily Updates Cron Output" cronreports.i9xgk@zapiermail.com'
0 2 * * * root /usr/bin/flock -n /tmp/update_frontend_xignite.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_xignite.php 2>&1 | tail -n 200 | mail -s "Xignite Daily Updates Cron Output" cronreports.i9xgk@zapiermail.com'
20 0,9,11,13,15,16 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_indexes.php 2>&1 | tail -n 200 | mail -s "Update Market Indexes" cronreports.i9xgk@zapiermail.com

##New external data process (EOL+GF)
0 7,20 * * * root /usr/bin/flock -n /tmp/eolgf_splits_checker.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/eolgf_splits_checker.php 2>&1 | tail -n 200 | mail -s "Split Checker Process" cronreports.i9xgk@zapiermail.com'
0 0 * * * root /usr/bin/flock -n /tmp/eolgf_updater.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/eolgf_updater.php 2>&1 | tail -n 200 | mail -s "Main Proedgar Updater Process" cronreports.i9xgk@zapiermail.com'
0 5 * * * root /usr/bin/flock -n /tmp/eolgf_updater_100days.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/eolgf_updater_100days.php 2>&1 | tail -n 200 | mail -s "Update Tickers 100 days old Process" cronreports.i9xgk@zapiermail.com'
0 8 * * * root /usr/bin/flock -n /tmp/eolgf_updater_375days.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/eolgf_updater_375days.php 2>&1 | tail -n 200 | mail -s "Update Tickers 375 days old Process" cronreports.i9xgk@zapiermail.com'
30 6,20 * * * root /usr/bin/flock -n /tmp/eolgf_updater_guru_miss.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/eolgf_updater_guru_miss.php 2>&1 | tail -n 200 | mail -s "Update Missing Guru reports Process" cronreports.i9xgk@zapiermail.com'

##Internal Data process
45 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_details.php
0 1 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_top.php 2>&1 | tail -n 200 | mail -s "Update Screener Top Filters Cron Output" cronreports.i9xgk@zapiermail.com
30 23 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_eod_valuation.php 2>&1 | tail -n 200 | mail -s "Update EOD Valuation Cron Output" cronreports.i9xgk@zapiermail.com
0 5 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_portfolio_dividends.php 2>&1 | tail -n 200 | mail -s "Update Portfolio Dividends Cron Output" cronreports.i9xgk@zapiermail.com
30 1 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/portfolio_cleanup.php 2>&1 | tail -n 200 | mail -s "Portfolio Cleanup" cronreports.i9xgk@zapiermail.com
45 7 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_industries.php 2>&1 | tail -n 200 | mail -s "Industry COC matching process Output" cronreports.i9xgk@zapiermail.com

##Other processes
12 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/check_yahoo_status.php
30 2 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /var/www/oldschoolvalue/middleware/crons/img_cleanup.sh 2>&1 | tail -n 200 | mail -s "Image Cache Cleanup Cron Output" cronreports.i9xgk@zapiermail.com
0 23 * * 5 root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/send_weekly_emails.php 2>&1 | tail -n 200 | mail -s "Send Weekly Emails" cronreports.i9xgk@zapiermail.com
