# crons/frontend_data.cron
0 6,18 * * * root /usr/bin/flock -n /tmp/update_frontend_yahoo_estimates.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_estimates.php 2>&1 | tail -n 200 | mail -s "Yahoo Estimates Cron Output" cronreports.i9xgk@zapiermail.com'
0 17,19,21,23 * * * root /usr/bin/flock -n /tmp/update_frontend_yahoo_daily.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_daily.php 2>&1 | tail -n 200 | mail -s "Yahoo Daily Updates Cron Output" cronreports.i9xgk@zapiermail.com'
30 4 * * * root /usr/bin/flock -n /tmp/update_frontend_EOL_GF_data.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_EOL_GF_data.php 2>&1 | tail -n 200 | mail -s "EOL/GF Daily Updates Cron Output" cronreports.i9xgk@zapiermail.com'
12 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/check_yahoo_status.php
0 2 * * * root /usr/bin/flock -n /tmp/update_frontend_xignite.lock sh -c 'cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_xignite.php 2>&1 | tail -n 200 | mail -s "Xignite Daily Updates Cron Output" cronreports.i9xgk@zapiermail.com'
45 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_details.php
0 1 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_top.php 2>&1 | tail -n 200 | mail -s "Update Screener Top Filters Cron Output" cronreports.i9xgk@zapiermail.com
30 23 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_eod_valuation.php 2>&1 | tail -n 200 | mail -s "Update EOD Valuation Cron Output" cronreports.i9xgk@zapiermail.com
30 2 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /var/www/oldschoolvalue/middleware/crons/img_cleanup.sh 2>&1 | tail -n 200 | mail -s "Image Cache Cleanup Cron Output" cronreports.i9xgk@zapiermail.com
30 3 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_portfolio_dividends.php 2>&1 | tail -n 200 | mail -s "Update Portfolio Dividends Cron Output" cronreports.i9xgk@zapiermail.com
20 0,9,11,13,15,16 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_indexes.php 2>&1 | tail -n 200 | mail -s "Update Market Indexes" cronreports.i9xgk@zapiermail.com
