# crons/frontend_data.cron
0 6,18 * * * root [ `ps -ef|grep -v grep | grep update_frontend_yahoo_estimates.php|wc -l` -eq 0 ] && cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_estimates.php 2>/dev/null
0 */4 * * * root [ `ps -ef|grep -v grep | grep update_frontend_yahoo_daily.php|wc -l` -eq 0 ] && cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_yahoo_daily.php 2>/dev/null
30 4 * * * root [ `ps -ef|grep -v grep | grep update_frontend_EOL_GF_data.php|wc -l` -eq 0 ] && cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_EOL_GF_data.php 2>/dev/null
12 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/check_yahoo_status.php
0 2 * * * root [ `ps -ef|grep -v grep | grep update_frontend_xignite.php|wc -l` -eq 0 ] && cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_frontend_xignite.php 2>/dev/null
45 * * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_details.php
0 1 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_screener_filters_top.php
0 22 * * * root cd /var/www/oldschoolvalue/middleware/crons/; /usr/bin/php /var/www/oldschoolvalue/middleware/crons/update_eod_valuation.php
